# mua-kids-chat/files/etc/init.d/mua
#!/bin/sh /etc/rc.commonSTART=99start(){    /usr/bin/mua_main_service 2>&1 > /dev/null &    sleep 1    /usr/bin/mua_mqtt_service 2>&1 > /dev/null &    /usr/bin/mua_miod_service 2>&1  > /dev/null &    /usr/bin/mua_file_service 2>&1 > /dev/null &    /usr/bin/mua_snd_service 2>&1 > /dev/null &    /usr/bin/mua_msg_service 2>&1 > /dev/null &    sleep 5    ubus call mua.main.service start {}}